# videocontrol.service
[Unit]
Description=VCServer
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=vcgroup
WorkingDirectory=/vid/videocontrol
ExecStart=/usr/bin/node /vid/videocontrol/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/vid/videocontrol/public/content /vid/videocontrol/config /vid/videocontrol/logs /vid/videocontrol/temp /vid/videocontrol/.converted /home/YOUR_USERNAME/.cache /home/YOUR_USERNAME/.config

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
