# clients/mpv/videocontrol-mpv@.service
[Unit]
Description=VideoControl MPV Client for %i
Documentation=https://github.com/ya-k0v/VideoControl
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=videocontrol
Group=videocontrol

# Переменные окружения
Environment="DISPLAY=:0"
Environment="PULSE_SERVER=/run/user/1000/pulse/native"

# Рабочая директория
WorkingDirectory=/opt/videocontrol-mpv

# Команда запуска
ExecStart=/usr/bin/python3 /opt/videocontrol-mpv/mpv_client.py \
    --server ${SERVER_URL} \
    --device %i \
    --display :0

# Автоматический перезапуск при крашах
Restart=always
RestartSec=5

# Systemd watchdog для мониторинга здоровья
# WatchdogSec=30

# Лимиты для предотвращения утечек памяти
MemoryMax=512M
CPUQuota=80%

# Автоматический restart при слишком частых падениях
StartLimitInterval=300
StartLimitBurst=5

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=videocontrol-mpv-%i

# Безопасность
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target

